# deployment/upstart/winchatty-frontend.conf
setuid USERNAME
start on runlevel [2345]
stop on runlevel [016]

chdir /home/chatty/frontend
exec npm start
respawn
