# deployment/upstart/winchatty-notify-server.conf
setuid USERNAME
start on runlevel [2345]
stop on runlevel [016]

chdir /mnt/websites/winchatty.com
exec /mnt/websites/winchatty.com/notify_server_run.sh
respawn
