# deployment/upstart/winchatty-search.conf
setuid USERNAME
start on runlevel [2345]
stop on runlevel [016]

chdir /home/chatty/search-data/
exec /usr/bin/dtsearchd 8081 a.dat b.dat c.dat
respawn
