# deployment/upstart/winchatty-push-server.conf
setuid USERNAME
start on runlevel [2345]
stop on runlevel [016]

exec /home/chatty/backend/push_run.sh
respawn
