# deployment/upstart/winchatty-push-server.conf
setuid USERNAME
start on runlevel [2345]
stop on runlevel [016]

chdir /mnt/websites/winchatty.com
exec /mnt/websites/winchatty.com/push_run.sh
respawn
