# deployment/upstart/winchatty-indexer.conf
setuid USERNAME
start on runlevel [2345]
stop on runlevel [016]

exec /home/chatty/backend/nusearch_run.sh
respawn
